# battery-monitor.service
[Unit]
Description=Battery Full Charger Alert Monitor
After=graphical.target sound.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/battery-monitor/battery_monitor.py
User=%USER%
Restart=on-failure
RestartSec=10
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native

[Install]
WantedBy=default.target
